# systemd/get_gps_data1.service
[Unit]
Description=Get GPS Location
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi

# GPS / MQTT config
Environment=GPS_DEBUG_LOG_RAW=0
Environment=GPS_FORCE_GPSD=0
Environment=MQTT_BROKER_HOST=api.copilotai.click
Environment=MQTT_PUBLISH_INTERVAL=2.0
Environment="DB_PATH=/home/pi/facial-tracker-firmware/data/blinksmart.db"

# OTA-safe paths
Environment=PYTHONPATH=/home/pi/facial-tracker-firmware/current/apps:/home/pi/facial-tracker-firmware/current/libs
WorkingDirectory=/home/pi/facial-tracker-firmware/current/apps

ExecStart=/home/pi/facial-tracker-firmware/venv/bin/python \
          /home/pi/facial-tracker-firmware/current/apps/get_gps_data.py

Restart=always
RestartSec=3

# Lower priority than facial tracking — needed for speed data but must not block I/O
Nice=5
CPUWeight=30
IOSchedulingClass=idle
IOSchedulingPriority=7
MemoryMax=30M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target