# systemd/ota-startup.service
[Unit]
Description=OTA Startup Service - Check for updates on boot
After=network-online.target
Wants=network-online.target
Before=facial1.service get_gps_data1.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=pi
Group=pi

WorkingDirectory=/home/pi/facial-tracker-firmware

ExecStart=/home/pi/facial-tracker-firmware/venv/bin/python \
          /home/pi/facial-tracker-firmware/ota_startup.py

StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target