# systemd/provisioning_ui.service
[Unit]
Description=Device Provisioning Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
Environment="PYTHONPATH=/home/pi/facial-tracker-firmware/src"
WorkingDirectory=/home/pi/facial-tracker-firmware/src

ExecStart=/home/pi/facial-tracker-firmware/venv/bin/python \
          /home/pi/facial-tracker-firmware/src/provisioning_ui.py

Restart=on-failure
RestartSec=5

# Lower priority than facial tracking (Nice=-5)
Nice=10
CPUWeight=25
MemoryMax=60M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target