# systemd/upload_images.service
[Unit]
Description=Upload Critical Images
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi

# OTA-safe paths
Environment=PYTHONPATH=/home/pi/facial-tracker-firmware/current/apps:/home/pi/facial-tracker-firmware/current/libs
WorkingDirectory=/home/pi/facial-tracker-firmware/current/apps

ExecStart=/home/pi/facial-tracker-firmware/venv/bin/python \
          /home/pi/facial-tracker-firmware/current/apps/upload_images.py

Restart=on-failure
RestartSec=5

# Lower priority than facial tracking (Nice=-5)
Nice=10
CPUWeight=25
MemoryMax=60M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target