# systemd/upload_images.service
[Unit]
Description=Upload Critical Images
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi

# OTA-safe paths
Environment="DB_PATH=/home/pi/facial-tracker-firmware/data/blinksmart.db"
Environment=PYTHONPATH=/home/pi/facial-tracker-firmware
WorkingDirectory=/home/pi/facial-tracker-firmware

ExecStart=/home/pi/facial-tracker-firmware/venv/bin/python \
          /home/pi/facial-tracker-firmware/upload_images.py

Restart=on-failure
RestartSec=5

# Lowest priority — bulk I/O must not compete with detection
Nice=15
CPUWeight=10
IOSchedulingClass=idle
IOSchedulingPriority=7
MemoryMax=25M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target